# rpm/systemd/lanecove-peer-1.service
[Unit]
Description=Lane Cove Tunnel - Peer 1
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/lanecove-peer -c /etc/lanecove/peer-1.yaml
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
